# deploy/gunicorn.service
[Unit]
Description=Mewada Matrimonial App – Gunicorn WSGI server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=deployuser
Group=deployuser
WorkingDirectory=/var/www/mewada
Environment="PATH=/var/www/mewada/venv/bin"
EnvironmentFile=/var/www/mewada/.env

ExecStart=/var/www/mewada/venv/bin/gunicorn \
    --workers 2 \
    --worker-class sync \
    --timeout 120 \
    --bind 127.0.0.1:5000 \
    --access-logfile /var/log/mewada/access.log \
    --error-logfile /var/log/mewada/error.log \
    "app:create_app()"

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
